# oraclenfs.mount
[Unit]
Description=Mount oraclenfs NFS volume to /oraclenfs

[Mount]
What=nfs:/oraclenfs
Where=/oraclenfs
Type=nfs
